# dotfiles/vscode/openclaw-code-tunnel.service
[Unit]
# VS Code Remote Tunnel for the GOTO 2026 student box. Students reach the editor
# at https://vscode.dev/tunnel/<VSCODE_TUNNEL_NAME> — GitHub OAuth secures the
# link, so no inbound port is exposed.
#
# The tunnel name is a 24-char hex string derived by clone.sh as
# sha256(hostname + TUNNEL_SALT + "vscode")[:24] — not the hostname, not in DNS,
# unguessable without both secrets. It is injected into
# /etc/openclaw/vscode-tunnel-name by cloud-init and read here via EnvironmentFile.
#
# Generic at bake: the unit is installed by phase_vscode but left DISABLED.
# cloud-init runcmd runs `code tunnel user login` then enables+starts this unit
# per box, so the golden image carries no token and no per-box state.
Description=VS Code Remote Tunnel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
EnvironmentFile=/etc/openclaw/vscode-tunnel-name
ExecStart=/usr/bin/code tunnel --accept-server-license-terms --name $VSCODE_TUNNEL_NAME
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
