# dotfiles/desktop/openclaw-desktop-cred.service
[Unit]
Description=Materialize classroom desktop basic-auth from per-instance credentials
# cloud-final writes /etc/openclaw/desktop.env; run after it, before nginx serves.
After=cloud-final.service
Before=nginx.service
Wants=nginx.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/openclaw-desktop-cred.sh

[Install]
WantedBy=multi-user.target
